# test.service
[Unit]
Description=Test unit for test timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test.sh