# test.timer
[Unit]
Description=A test timer

[Timer]
OnCalendar=*-*-* *:0/5

[Install]
WantedBy=timers.target